# deploy/ansible/roles/queue_worker/templates/worker.service.j2
[Unit]
Description=KewlChats queue worker ({{ app_name }})
After=network.target redis-server.service mysql.service ejabberd.service
Requires=redis-server.service

[Service]
User={{ app_name }}
Group={{ app_name }}
Restart=always
RestartSec=5
WorkingDirectory={{ app_root }}/{{ app_name }}
ExecStart=/usr/bin/php{{ php_version }} {{ app_root }}/{{ app_name }}/artisan queue:work redis \
  --tries=5 --backoff=10 --timeout=120 --sleep=3 --max-jobs=1000 --max-time=3600
# Graceful restart: let the worker finish the current job.
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target
