# deploy/ansible/roles/ejabberd/templates/kewlchats-mam-trim.timer.j2
[Unit]
Description=Daily MAM trim (kewlchats)

[Timer]
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
