# deploy/ansible/roles/certbot/templates/kewlchats-certbot-renew.timer.j2
[Unit]
Description=Twice-daily Let's Encrypt renewal check (kewlchats)

[Timer]
OnCalendar=*-*-* 03,15:00:00
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
